# root/etc/systemd/user/wireplumber.service.d/override.conf
[Service]
Environment="DISABLE_RTKIT=1"
