# root/etc/systemd/user/session-manager.service
[Unit]
Description=Startup service

[Service]
ExecStart=/usr/local/bin/session-manager

[Install]
WantedBy=default.target
