# root/etc/systemd/user/xvfb.service
[Unit]
Description=Run Xvfb

[Service]
EnvironmentFile=/etc/gnome-docker/env
ExecStart=/usr/local/bin/xvfb

[Install]
WantedBy=default.target
