# etc/systemd/system/systemd-logind.service.d/override.conf
# Systemd requires this for work without CAP_SYS_ADMIN

[Service]
ProtectHostname=no
